# linux/rootfs/etc/inittab
id:3:initdefault：
#console:2345:respawn:/bin/login console
